# filesystem/etc/update-motd.d/10-uname
#!/bin/bash
echo -e "\e[1;36m$(uname -srv)\e[0m"
